# package/mistify/etcd/etcd-setup.service
[Unit]
Description=Etcd Initialization and Setup

[Service]
ExecStart=/usr/lib/systemd/scripts/etcd-setup

[Install]
WantedBy=default.target
