# board/mistify/rootfs_overlay/root/lochness-ansible/roles/etcd/templates/etcd-announce.service.j2
[Unit]
Description=Announce etcd
BindsTo=etcd.service

[Service]
ExecStart=/usr/sbin/queensland announce etcd-server --name={{machine_id}} --port=2380
User=nobody
RestartSec=30s
Restart=always
