# board/mistify/rootfs_overlay/root/lochness-ansible/roles/lochness/templates/node-announce.service.j2
[Unit]
Description=Announce Node
Wants=etcd.service

[Service]
ExecStart=/usr/sbin/queensland node --name={{machine_id}} --address={{ ansible_default_ipv4.address }}
User=nobody
RestartSec=30s
Restart=always
