# board/mistify/rootfs_overlay/etc/inittab
# /etc/inittab
#
# Copyright (C) 2001 Erik Andersen <andersen@codepoet.org>
#
# Note: BusyBox init doesn't support runlevels.  The runlevels field is
# completely ignored by BusyBox init. If you want runlevels, use
# sysvinit.
#
# Format for each entry: <id>:<runlevels>:<action>:<process>
#
# id        == tty to run on, or empty for /dev/console
# runlevels == ignored
# action    == one of sysinit, respawn, askfirst, wait, and once
# process   == program to run
id:1:initdefault:

# Startup the system
proc::sysinit:/bin/mount -t proc proc /proc
rwmo::sysinit:/bin/mount -o remount,rw / # REMOUNT_ROOTFS_RW
dpts::sysinit:/bin/mkdir -p /dev/pts
dshm::sysinit:/bin/mkdir -p /dev/shm
moun::sysinit:/bin/mount -a
host::sysinit:/bin/hostname -F /etc/hostname

# now run any rc startup scripts
init::sysinit:/etc/init.d/rcS

# Local TTYs
1:1:respawn:/sbin/getty 38400 tty1
2:1:respawn:/sbin/getty 38400 tty2

# Put a getty on the serial port
S1::respawn:/sbin/getty -L ttyS0 115200 vt100 # GENERIC_SERIAL

# Logging stuff
mess::sysinit:/bin/touch /var/log/mesages

# Stuff to do for the 3-finger salute
rebo::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
hlt:0:wait:/etc/init.d/rcK-halt
rb:6:wait:/etc/init.d/rcK-reboot
