# board/mistify/rootfs_overlay/root/lochness-ansible/roles/agent/templates/mistify-agent-image.service.j2
[Unit]
Description=Mistify Image Sub-agent
After=network.target zfs.target

[Service]
Type=simple
ExecStart=/opt/mistify/sbin/mistify-agent-image
LimitNOFILE=32768
RestartSec=60s
Restart=always

[Install]
WantedBy=multi-user.target
