# package/mistify/confd/confd-setup.service
[Unit]
Description=ConfD Initialization and Setup

[Service]
ExecStart=/usr/lib/systemd/scripts/confd-setup

[Install]
WantedBy=default.target
