# board/mistify/rootfs_overlay/etc/systemd/system/vtun.service
[Unit]
Description=virtual tunnel daemon
After=syslog.target

[Service]
Type=forking
ExecStart=/usr/sbin/vtund -s

[Install]
WantedBy=multi-user.target
